# package2/supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /prod/project2/wsgi.ini --logto /tmp/uwsgi.log

[program:nginx-app]
command = /usr/sbin/nginx